# tests/deb/supervisord.conf
[supervisord]
nodaemon=true

[program:install_tao]
autorestart=false
directory=/usr/share/tao-community
command=/bin/bash -c 'until mysql  -uroot -ptao -htao-db -e"SELECT 1 FROM DUAL"; do echo Waiting tao-db to start; sleep 2; done;  /usr/bin/php tao/scripts/taoInstall.php --db_driver pdo_mysql --db_host tao-db --db_name tao --db_user root --db_pass tao -e taoCe -u admin -p admin -url http://localhost:8740/tao-community/ && echo && echo -e "\\n\\nTAO installed successfully.\\nGo to http://localhost:8740/tao-community/ and login as admin/admin\\n" 2>&1'
user=taosrv
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:sshd]
command=/usr/sbin/sshd -D

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:php-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize --fpm-config /etc/php/7.0/fpm/php-fpm.conf
